# pi/com.concept2.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.concept2.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/bjallen/Sites/concept2-mqtt/.venv/bin/python3</string>
        <string>/Users/bjallen/Sites/concept2-mqtt/monitor.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/bjallen/Sites/concept2-mqtt</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/bjallen/Sites/concept2-mqtt/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bjallen/Sites/concept2-mqtt/monitor.log</string>
</dict>
</plist>
